# doc/crontab
@reboot /usr/bin/supervise ~/abtf.com > /dev/null 2>&1
