# node/crontab
@reboot /usr/bin/supervise ~/abtf.com/node > /dev/null 2>&1
